# etc/rc.local
#!/bin/sh

/home/ubuntu/RDFeasy/sbin/first_boot_initialization
exit 0
