# etc/rc.local
#!/bin/sh

/home/ubuntu/RDFeasy/sbin/change_password_at_boot

exit 0
